# com.meetingmind.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!--
        MeetingMind Launch Agent
        Starts the daemon automatically on user login.

        Installation:
          cp com.meetingmind.agent.plist ~/Library/LaunchAgents/
          launchctl load ~/Library/LaunchAgents/com.meetingmind.agent.plist

        Removal:
          launchctl unload ~/Library/LaunchAgents/com.meetingmind.agent.plist
          rm ~/Library/LaunchAgents/com.meetingmind.agent.plist
    -->
    <key>Label</key>
    <string>com.meetingmind.agent</string>

    <key>ProgramArguments</key>
    <array>
        <!--
            Production (PyInstaller binary):
              <string>/Applications/MeetingMind.app/Contents/Resources/meetingmind-daemon</string>

            Development (Python venv):
              <string>/path/to/meeting-mind/.venv/bin/python3</string>
              <string>-m</string>
              <string>src.main</string>
        -->
        <string>/Applications/MeetingMind.app/Contents/Resources/meetingmind-daemon</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Applications/MeetingMind.app/Contents/Resources</string>

    <!-- Start automatically on login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart if the process exits unexpectedly. -->
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <!-- Throttle restarts to avoid rapid-fire loops on error. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- Log to ~/Library/Logs for easier debugging. -->
    <key>StandardOutPath</key>
    <string>/Users/USER/Library/Logs/meetingmind/stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/USER/Library/Logs/meetingmind/stderr.log</string>
</dict>
</plist>
